# systemd/terraform-webhook.service
[Unit]
Description=Terraform Runner Webhook Server
After=network.target

[Service]
Type=simple
User=terraform
Group=terraform
WorkingDirectory=/opt/terraform-runner
ExecStart=/usr/bin/python3 /opt/terraform-runner/webhook-server.py
Restart=always
RestartSec=10

# Environment variables
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/etc/environment

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=terraform-runner-webhook

[Install]
WantedBy=multi-user.target